# roles/octest/files/octest.conf

description "octest"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid ubuntu

exec /home/ubuntu/scripts/octest
